# container-files/services.d/syslog-ng/run
#!/bin/execlineb -P

/usr/sbin/syslog-ng --cfgfile=/etc/syslog-ng/syslog-ng.conf -F